# roles/snapshots/files/snapshot-timeline.timer

[Unit]
Description=Timeline of snapshots for /home/

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target

