# roles/snapshots/files/snapshot-timeline.service

[Unit]
Description=Timeline of snapshots for /home/

[Service]
Type=simple
ExecStart=/usr/lib/snapshot-timeline

