# roles/restic-backup/templates/restic-backup.timer.j2
[Unit]
Description=Backup user home to NFS
RefuseManualStart=no
RefuseManualStop=no

[Timer]
OnBootSec=5m
OnUnitActiveSec=4h
Unit=restic-backup.service

[Install]
WantedBy=timers.target