# roles/restic-backup/templates/restic-backup.service.j2
[Unit]
Description=Backup user home to NFS

[Service]
Type=simple
ExecStart=/home/{{ user.name }}/apps/backup/backup.bash

[Install]
WantedBy=default.target